# lib/linux/servers/queries/palantir/servers/crontab.yml
- name: crontab
  description: Retrieves all the jobs scheduled in crontab in the target system.
  query: SELECT * FROM crontab;
  platform: linux
