# lib/macos/both/queries/osquery-configuration/endpoints/launchd.yml
- name: launchd
  description: LaunchAgents and LaunchDaemons from default search paths.
  query: SELECT * FROM launchd;
  platform: darwin
