# lib/macos/both/queries/osquery-configuration/endpoints/startupitems.yml
- name: startup_items
  description: Applications and binaries set as user/login startup items.
  query: SELECT * FROM startup_items;
  platform: darwin
