# lib/all/queries/palantir-configuration/informational/startupitems.yml
---
apiVersion: v1
kind: query
spec:
  name: startup_items
  description: Applications and binaries set as user/login startup items.
  query: SELECT * FROM startup_items;
  team: both
