# lib/macos/queries/osquery-configuration/endpoints/launchd.yml
---
apiVersion: v1
kind: query
spec:
  name: launchd
  description: LaunchAgents and LaunchDaemons from default search paths.
  query: SELECT * FROM launchd;
  platform: darwin
